# cfenginefix/check_cfengine.cron
0 */8 * * * root <path>/fix_cfengine.sh
